# .config/systemd/user/autotheme.service
[Unit]
Description=Auto theme switch

[Service]
Type=oneshot
ExecStart=%h/.local/bin/autotheme
