# .config/systemd/user/darkmode.service
[Unit]
Description=Auto dark mode switch

[Service]
Type=oneshot
ExecStart=%h/.local/bin/darkmode auto
