# .config/systemd/user/theme.timer
[Unit]
Description=Theme setup periodically

[Timer]
OnStartupSec=1m
OnUnitActiveSec=1m
Persistent=true

[Install]
WantedBy=timers.target
