# .config/systemd/user/theme.service
[Unit]
Description=Theme setup

[Service]
Type=oneshot
ExecStart=/usr/bin/systemd-run --scope --unit=themetmp %h/.local/bin/theme setup
