# .config/systemd/user/theme.service
[Unit]
Description=Theme setup

[Service]
Type=oneshot
ExecStart=%h/.local/bin/theme setup
PassEnvironment=DISPLAY WAYLAND_DISPLAY XDG_RUNTIME_DIR
