# .config/systemd/user/themeauto.service
[Unit]
Description=Auto switch theme

[Service]
Type=oneshot
ExecStart=%h/.local/bin/theme auto
