# .config/systemd/user/theme.service
[Unit]
Description = Theme setup

[Service]
Type = oneshot
KillMode = none
ExecStart = %h/.local/bin/theme setup
