# .config/systemd/user/theme.service
[Unit]
Description=Theme setup

[Service]
Type=oneshot
KillMode=process
ExecStart=%h/.local/bin/theme setup
