# .config/systemd/user/wallpaper.service
[Unit]
Description=Change wallpaper

[Service]
Type=oneshot
ExecStart=%h/.local/bin/wallpaper
