# .config/systemd/system/getty@.service.d/skip-username.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty -o '-p -- yahngming' --noclear --skip-login - $TERM
